# roles/llamaindex-memory-worker/templates/llamaindex-memory-worker.service.j2
[Unit]
Description=LlamaIndex memory worker
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ memory_worker_user }}
Group={{ memory_worker_user }}
WorkingDirectory={{ memory_worker_install_dir }}
EnvironmentFile={{ memory_worker_env_file }}
# Block until Waza loadavg is below memory_worker_loadavg_threshold.
# If it stays hot, ExecStartPre fails and the run is skipped. The
# memory-healthcheck n8n workflow detects the resulting stale report.
ExecStartPre={{ memory_worker_wait_calm_script }}
ExecStart={{ memory_worker_run_report_script }} --mode {{ memory_worker_index_mode }}
Nice=19
IOSchedulingClass=idle

# Resource control (cgroup) — REX 2026-06-05 panne SSH waza (parité avec l'unit
# user memory-worker.service.user.j2 — caps D14, contrat RAG v3 2026-06-10).
MemoryHigh={{ memory_worker_memory_high }}
MemoryMax={{ memory_worker_memory_max }}
MemorySwapMax={{ memory_worker_memory_swap_max }}
OOMScoreAdjust={{ memory_worker_oom_score_adjust }}

[Install]
WantedBy=multi-user.target
