# roles/mission-control/templates/mission-control.service.j2
[Unit]
Description=Mission Control — OpenClaw Web UI
After=network.target

[Service]
Type=simple
User={{ mc_user }}
WorkingDirectory={{ mc_install_dir }}
EnvironmentFile={{ mc_install_dir }}/.env
ExecStart=/usr/bin/node {{ mc_install_dir }}/node_modules/.bin/next start -p {{ mc_port }}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mission-control

[Install]
WantedBy=multi-user.target
